# packaging/windows_scheduler_template.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  StocksTrading 每日例行任務範本 (Windows Task Scheduler).

  匯入方式：
    1. 編輯下方 <Author>、<UserId>、<Command>、<WorkingDirectory>、<StartBoundary>
       為你的環境設定 (見 packaging/README.md)
    2. 開啟「工作排程器」(taskschd.msc)
    3. 動作 → 匯入工作... → 選此 XML
    4. 觸發程序頁籤確認時間 (預設每天 18:00 跑)

  注意：
    - StartBoundary 是台股收盤 + 安全餘裕的時間；美股市場請改為 06:00 隔日
    - 預設 dry-run，要實際寫入訊號 / 寄信，移除 --dry-run 即可
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Author>YOUR_USERNAME</Author>
    <Description>每日 18:00 跑 stocks-trading-cli daily-routine，產生交易訊號並寄送日報</Description>
    <URI>\StocksTrading\DailyRoutine</URI>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <!-- 第一次啟動時間；之後依 ScheduleByDay 每日重複 -->
      <StartBoundary>2026-06-01T18:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>YOUR_USERNAME</UserId>
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT15M</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <!-- 請改成你安裝後的 stocks-trading-cli.exe 路徑 -->
      <Command>C:\Program Files\StocksTrading\stocks-trading-cli.exe</Command>
      <Arguments>daily-routine --tickers SPY,QQQ,IWM --strategy dual-momentum --dry-run</Arguments>
      <WorkingDirectory>C:\Program Files\StocksTrading</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
